# CR-PP-GPL-20221022/src-comcerto-openwrt-ls1024_7.0.0-rc0-gpl/package/freescale/c2k_gpio_reset/files/etc/init.d/gpio_reset
#!/bin/sh /etc/rc.common
# 2017-02-06: harry.lin <harry.lin@deltaww.com>
# Add c2k_gpio_reset module support.
# Copyright (C) 2008 OpenWrt.org

START=10

start() {
        insmod c2krv16x_26x_gpio_reset
}

stop(){
        rmmod c2krv16x_26x_gpio_reset
}

